# dpkg-skel/ccscanner_VERSION-0ubuntu_arm64/lib/systemd/system/cc-nmap-taskmanager.service
[Unit]
Description=command center nmap task-manager service
After=network.target

[Service]
Restart=always
RestartSec=1
User=ccscanner
ExecStart=/usr/bin/cc-nmap-taskmanager --config=/etc/ccscanner/config.json

[Install]
WantedBy=multi-user.target

